# rootfs/etc/rc.local
#!/bin/bash

PATH=/nix/var/nix/profiles/default/bin/:$PATH

# Turn off alt+arrows to switch VTs, and shift+nav keys to scroll
dumpkeys -k \
| egrep -v '(Incr|Decr|Last)_Console' \
| egrep -v 'Scroll_(Backward|Forward)' \
| egrep -v KeyboardSignal \
| loadkeys

momovt config \
  width 5 height 3 \
  wrap 0 \
  splash_ttl 1 \
  splash_info momovt-banner

#systemctl start getty@tty{1..16}
